# cfg/demo-aq-main.conf
[program:demo-aq-main]
command=/usr/bin/php8.2 /root/demo/air_quality/artisan serve --host=0.0.0.0 --port=8005
numprocs=1
autostart=true
autorestart=true
user=root
