# cfg/demo-aq-broadcast-data.conf
[program:demo-aq-broadcast-data]
command=/usr/bin/php8.2 /root/demo/air_quality/artisan app:broadcast-data
numprocs=1
autostart=true
autorestart=true
user=root
