# cfg/demo-aq-broadcast-chart.conf
[program:demo-aq-broadcast-chart]
command=/usr/bin/php8.2 /root/demo/air_quality/artisan app:sensor-update
numprocs=1
autostart=true
autorestart=true
user=root
