# cfg/demo-aq-dummy-data.conf
[program:demo-aq-dummy-data]
command=/usr/bin/php8.2 /root/demo/air_quality/artisan app:dummy-data
numprocs=1
autostart=true
autorestart=true
user=root
